# crontab
0 22 * * * osascript /Users/lvarayut/Documents/evernote-backup/EvernoteBackup.scpt
